# tools/book-image-service/com.coderdream.book-image-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.coderdream.book-image-service</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/System/docker/book-image-service/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>app:app</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>30019</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Volumes/System/docker/book-image-service</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Volumes/System/docker/book-image-service/service.log</string>
  <key>StandardErrorPath</key>
  <string>/Volumes/System/docker/book-image-service/service.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HF_HUB_OFFLINE</key>
    <string>1</string>
    <key>HF_ENDPOINT</key>
    <string>https://hf-mirror.com</string>
    <key>BOOK_IMAGE_DTYPE</key>
    <string>auto</string>
    <key>BOOK_IMAGE_MODEL</key>
    <string>Lykon/dreamshaper-8-lcm</string>
    <key>BOOK_IMAGE_OUTPUT_DIR</key>
    <string>/Volumes/System/docker/book-image-service/outputs</string>
  </dict>
</dict>
</plist>
